# scrypted-beta/rootfs/etc/services.d/scrypted-run/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Log Viewer
# Runs the Log Viewer
# ==============================================================================
declare name
declare value
bashio::log.info "Starting Scrypted..."

# Ensure the configuration exists
if bashio::fs.file_exists '/config/data_scrypted/options.json'; then
    bashio::log.info "Copying options.json"
    cp -f /config/data_scrypted/options.json /etc/options
else
    bashio::log.info "Making config directory"
    mkdir -p -v /config/data_scrypted \
        || bashio::exit.nok "Failed to create the Scrypted configuration directory"

    # Copy in template file
    bashio::log.info "Copying template file"
    cp -v /etc/options.json /config/data_scrypted/
fi


for var in $(bashio::config 'env_vars|keys'); do
    name=$(bashio::config "env_vars[${var}].name")
    value=$(bashio::config "env_vars[${var}].value")
    bashio::log.info "Setting ${name} in scrypted run" 
    export "${name}=${value}"
done

cd /server && npm exec scrypted-serve
