# scrypted-ingress-beta/rootfs/etc/services.d/scrypted-run/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Log Viewer
# Runs the Log Viewer
# ==============================================================================
bashio::log.info "Starting Scrypted..."

# Ensure the configuration exists
if bashio::fs.file_exists '/config/data_scrypted/options.json'; then
    cp -f /config/data_scrypted/options.json /etc/options
else
    mkdir -p /config/data_scrypted \
        || bashio::exit.nok "Failed to create the Scrypted configuration directory"

    # Copy in template file
    cp /etc/options.json /config/data_scrypted/
fi

cd /server && npm exec scrypted-serve
